# supervisord.conf
[supervisorctl]
;serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket 
serverurl=http://0.0.0.0:9001 ;
[inet_http_server]
port = 0.0.0.0:9001
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

# [program:danted]
# command=danted -f /etc/danted.conf
# autostart=true
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true

[program:squid]
command=/squid_entrypoint.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:openvpn]
command=openvpn --config /vpn.ovpn
;command=/bin/sh -c "mkdir -p /dev/net && mknod /dev/net/tun c 10 200 && exec openvpn --config /vpn.ovpn"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
